# config/supervisor/process.conf
[program:process]
command=php /var/www/bender/bin/console app:process
user=root
numprocs=1
autostart=true
autorestart=true
process_name=bender_%(program_name)s_%(process_num)02d